# deploy/com.quota-butler.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd 模板 · 每 __INTERVAL__ 秒拉起一次 quota_butler.main，跑完即退。
  占位符由 deploy/install.sh 在安装时填充，不要手动改这个文件。
  ⚠️ EnvironmentVariables 必填两项（本机实测，env -i 复现）：
     - PATH：launchd 默认 PATH 不含 /usr/local/bin，否则 main.py 调 lark-cli / claude
       会 FileNotFoundError。
     - LARK_CHANNEL：lark-cli 靠它选中"在群里的 bridge bot"；缺了会回退到默认 app，
       报 230002 "Bot can NOT be out of the chat."（推送静默失败）。
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.quota-butler</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>quota_butler.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
        <key>LARK_CHANNEL</key>
        <string>__LARK_CHANNEL__</string>
    </dict>

    <key>StartInterval</key>
    <integer>__INTERVAL__</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__REPO__/quota-butler.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO__/quota-butler.err.log</string>
</dict>
</plist>
